# charge-log.openrc
#!/sbin/openrc-run

description='Battery charge history logger'

command='/usr/bin/charge-log'
command_args='daemon'
command_background='true'
pidfile="run/${RC_SVCNAME}.pid"
