# charge-log.service.in
[Unit]
Description=Battery Charge History Logger

[Service]
ExecStart=%BINDIR%/charge-log-daemon
Restart=always

[Install]
WantedBy=multi-user.target
