# charge-log.openrc.in
#!/sbin/openrc-run

description='Battery charge history logger'

command='%BINDIR%/charge-log-daemon'
command_background='true'
pidfile="%RUNSTATEDIR%/${RC_SVCNAME}.pid"
